# root/etc/services.d/sshd/run
#!/bin/bash
set -x

chown -R app:app /config
usermod -d /config -s /bin/bash app

exec /usr/sbin/sshd -D -e
